# infrastructure/pocwisper.service
# Pocwisper systemd service file for AlmaLinux/RHEL
# This allows Pocwisper to run as a system service

[Unit]
Description=Pocwisper Audio Transcription Service
After=network.target podman.socket
Requires=podman.socket

[Service]
Type=forking
RemainAfterExit=yes
WorkingDirectory=/opt/pocwisper
ExecStart=/usr/bin/podman-compose -f /opt/pocwisper/podman-compose.yml up -d
ExecStop=/usr/bin/podman-compose -f /opt/pocwisper/podman-compose.yml down
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=multi-user.target
